# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

# Garage server + one-time self-bootstrap (layout / bucket / key).
[program:garage]
command=/server-init.sh
autostart=true
autorestart=unexpected
startsecs=10
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# Admin / management / browse web UI. Talks to the local admin + S3 APIs.
# Env (with defaults) is applied by the wrapper to avoid supervisor's brittle
# %(ENV_*)s expansion, which crashes when a variable is unset.
[program:garage-webui]
command=/webui-run.sh
autostart=true
autorestart=true
startretries=10
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
